# vagrant/puppet/modules/twemproxy/templates/nutcracker.erb
#!/bin/sh
### BEGIN INIT INFO
# Provides:          <%= name %>
# Required-Start:    $network $remote_fs $local_fs
# Required-Stop:     $network $remote_fs $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Stop/start nutcracker <%= name %>
### END INIT INFO

PATH=/sbin:/usr/sbin:/bin:/usr/bin
DESC=<%= name %>
NAME=<%= name %>
CONFFILE=/etc/nutcracker/<%= name %>.yml
LOGFILE=/var/log/nutcracker/<%= name %>.log
DAEMON=/usr/local/bin/nutcracker
PIDFILE=/var/run/nutcracker/$NAME.pid
DAEMON_ARGS="-c $CONFFILE -o $LOGFILE -p $PIDFILE -s <%= statsport %> -d"
#DAEMON_ARGS="-c $CONFFILE -p $PIDFILE -s <%= statsport %> -d"
SCRIPTNAME=/etc/init.d/$NAME

[ -x $DAEMON ] || exit 0

[ -r /etc/default/$NAME ] && . /etc/default/$NAME

. /lib/init/vars.sh

. /lib/lsb/init-functions

do_start()
{
    start-stop-daemon --start --quiet --pidfile $PIDFILE --exec $DAEMON -- \
        $DAEMON_ARGS
    RETVAL="$?"
    return "$RETVAL"
}

do_stop()
{
    # Return
    #   0 if daemon has been stopped
    #   1 if daemon was already stopped
    #   2 if daemon could not be stopped
    #   other if a failure occurred

    kill -9 `cat ${PIDFILE}`
    RETVAL="$?"
    rm -f ${PIDFILE}
    return "$RETVAL"
}

case "$1" in
    start)
        [ "$VERBOSE" != no ] && log_daemon_msg "Starting $DESC " "$NAME"
        do_start
        case "$?" in
            0|1) [ "$VERBOSE" != no ] && log_end_msg 0 ;;
            2) [ "$VERBOSE" != no ] && log_end_msg 1 ;;
        esac
        ;;
    stop)
        [ "$VERBOSE" != no ] && log_daemon_msg "Stopping $DESC" "$NAME"
        do_stop
        case "$?" in
            0|1) [ "$VERBOSE" != no ] && log_end_msg 0 ;;
            2) [ "$VERBOSE" != no ] && log_end_msg 1 ;;
        esac
        ;;
    status)
        ps -ef | grep `cat ${PIDFILE}`
        case "$?" in
            0|1)
                echo "Twemproxy ${NAME} is running..."
            ;;
            2)
                echo "Twemproxy ${NAME} isn't running..."
            ;;
         esac
        ;;
    restart)
        log_daemon_msg "Restarting $DESC" "$NAME"
        do_stop
        echo "Twemproxy has stopped..."
        do_start
                case "$?" in
                    0) log_end_msg 0 ;;
                    1) log_end_msg 1 ;; # Old process is still running
                    *) log_end_msg 1 ;; # Failed to start
                esac
    ;;
    *)
        echo "Usage: $SCRIPTNAME {start|stop|status|restart}" >&2
        exit 3
        ;;
esac

exit
$RETVAL
